# docs/crontab
* * * * * cd /tmp/test && /tmp/test/bin/shepherd launcher <api-ip>:<api-port>
